# ansible/roles/gsender/templates/gsender.service.j2
[Unit]
Description=gSender CNC Control Server
Documentation=https://github.com/Sienci-Labs/gsender
After=network.target

[Service]
Type=simple
User={{ gsender_user }}
Group={{ gsender_group }}
WorkingDirectory={{ gsender_home }}
ExecStart=/usr/bin/xvfb-run -a /opt/gSender/gsender \
    --host {{ gsender_host }} \
    --port {{ gsender_port }} \
    --remote \
{% if gsender_allow_remote_access %}
    --allow-remote-access \
{% endif %}
{% if gsender_controller %}
    --controller {{ gsender_controller }} \
{% endif %}
    --config {{ gsender_config_file }}

# Restart policy
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ gsender_home }}
SupplementaryGroups=dialout

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gsender

[Install]
WantedBy=multi-user.target
